# systemd/mach6-unified.service
[Unit]
Description=Symbiote Unified Daemon (Symbiote 3.0)
After=network.target

[Service]
Type=simple
User=adam
WorkingDirectory=/home/adam/workspace/symbiote/symbiote-core
ExecStart=/usr/bin/node dist/daemon/unified-daemon.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
